# .config/systemd/user/aria2.service
[Unit]
Description=Aria2 Daemon

[Service]
ExecStart=/usr/bin/aria2c

[Install]
WantedBy=default.target
